# recipes-core/base-files/files/systemd-hostnamed-watcher.path
[Unit]
After=data.mount
Requires=data.mount

[Path]
PathChanged=@@MENDER/CONST_CONF_HOSTNAME_CONF_DATA_DIR@@
MakeDirectory=True

[Install]
WantedBy=multi-user.target
