# recipes-containers/docker-load-archived/files/docker-load-archived.path
[Unit]
After=data.mount
Requires=data.mount

[Path]
PathChanged=@@MENDER/CONST_CONF_DOCKER_IMG_DATA_DIR@@
MakeDirectory=True

[Install]
WantedBy=multi-user.target
