# recipes-core/systemd-conf/files/systemd-networkd-watcher.service
[Unit]
Description=update systemd-networkd from /data
After=data.mount
Requires=data.mount systemd-networkd.service

[Service]
Type=oneshot
ExecStartPre=-@@base_bindir@@/mount -o remount @@sysconfdir@@/systemd/network
ExecStart=@@base_bindir@@/systemctl restart systemd-networkd.service

[Install]
WantedBy=multi-user.target
