# recipes-core/systemd-conf/files/systemd-networkd-watcher.path
[Path]
PathModified=@@MENDER/CONST_CONF_NETWORK_CONF_DATA_DIR@@
MakeDirectory=True

[Install]
WantedBy=multi-user.target
